# server/etc/cron.d/stunnel_presharedkey
# refresh pre shared keys for stunnel

MAILTO=root

45 03 * * *  root /usr/bin/rmapctrl --exportmqttpsk > /etc/stunnel/file.psk.new; cmp /etc/stunnel/file.psk /etc/stunnel/file.psk.new; if [ $? -ne 0 ]; then echo "stunnel restart"; mv -f /etc/stunnel/file.psk.new /etc/stunnel/file.psk; /bin/systemctl restart stunnel.service; fi
