# server/etc/systemd/postgresql.service.d/rmap.conf
[Service]
 
 # Location of database directory
 Environment=PGDATA=/rmap/pgsql/data
